# services.d/gopass-agent/finish
#!/command/execlineb -P
importas -i XDG_RUNTIME_DIR XDG_RUNTIME_DIR
rm -f ${XDG_RUNTIME_DIR}/gopass/gopass-age-agent.sock
